# woof-code/rootfs-packages/mac_changer/usr/lib/udev/rules.d/81-change-mac.rules
ACTION=="add", SUBSYSTEM=="net", RUN+="/usr/local/sbin/change-mac"
