# woof-code/rootfs-packages/swaybg_wallpaper/root/.config/autostart/swaybg.desktop
[Desktop Entry]
Version=1.0
Name=swaybg
Comment=swaybg
Exec=swaybg -c "#112222"
Terminal=false
Type=Application
