# woof-code/rootfs-skeleton/root/.config/autostart/swayidle.desktop
[Desktop Entry]
Version=1.0
Name=swayidle
Comment=swayidle
Exec=swayidle
Terminal=false
Type=Application
